# bot.service
[Unit]
Description=VK BOT
After=network.target

[Service]
User=root

WorkingDirectory=/root/Documents/vkbot/bot_main.py
Environment="PYTHONPATH=/root/Documents/vkbot/bot_main.py"
ExecStart=/usr/bin/python3 /root/Documents/vkbot/bot_main.py
 
[Install]
WantedBy=multi-user.target
