# bot.service
[Unit]
Description=VK BOT
After=network.target

[Service]
User=root
Group=root

WorkingDirectory=/usr/local/bin/bot/vkbot/
Environment="PYTHONPATH=/usr/local/bin/bot/vkbot/bot_main.py"
ExecStart=/usr/bin/python3.8 /usr/local/bin/bot/vkbot/bot_main.py
Restart=always
RestartSec=15
 
[Install]
Type=oneshot
WantedBy=multi-user.target
